# system/backupninja.timer
[Unit]
Description = Timer for backupninja

[Timer]
OnCalendar = daily

[Install]
WantedBy = multi-user.target
