# system/cpufreqd.service
[Unit]
Description = Userspace control of CPU power management

[Service]
ExecStart = /usr/sbin/cpufreqd -D

[Install]
WantedBy = multi-user.target

