# system/cron-daily.timer
[Unit]
Description = Daily cron job timer

[Timer]
OnCalendar = *-*-* 05:12:00

[Install]
WantedBy = multi-user.target
