# system/nslcd.service
[Unit]
Description = Name Service LDAP connection daemon

[Service]
ExecStart = /usr/sbin/nslcd -n

[Install]
WantedBy = multi-user.target

