# system/radicale.service
[Unit]
Description = CalDAV and CardDAV server

[Service]
ExecStart = /usr/bin/radicale
User = radicale

[Install]
WantedBy = multi-user.target
