# system/repro.service
[Unit]
Description = SIP proxy server
After = ifup@sit0.service

[Service]
ExecStart = /usr/sbin/repro

[Install]
WantedBy = multi-user.target

