# system/repro.service
[Unit]
Description = SIP proxy server

[Service]
ExecStart = /usr/sbin/repro /etc/repro/repro.config

[Install]
WantedBy = multi-user.target

