# system/saned@.service
[Unit]
Description = SANE remote scanner access
Requires = saned.socket

[Service]
ExecStart=/usr/sbin/saned
StandardInput = socket

